# rootfs/etc/s6-overlay/s6-rc.d/init-services/run
#!/command/with-contenv bashio
# shellcheck shell=bash

ln -s /data/services.d/* /run/service/ || true

kill -HUP 1
