# rootfs/etc/s6-overlay/s6-rc.d/init-services/run
#!/command/with-contenv bashio
# shellcheck shell=bash

SERVICES=/data/services.d

[ -d $SERVICES -a "$(ls -A $SERVICES/)" ] && ln -s $SERVICES/* /run/service/ || true
kill -HUP 1
